# agent/install/com.vant.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- VANT Signage Agent — macOS launchd plist
     Install:
       sudo cp com.vant.agent.plist /Library/LaunchDaemons/
       sudo launchctl load /Library/LaunchDaemons/com.vant.agent.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vant.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/vant-agent</string>
        <string>run</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VANT_CONFIG</key>
        <string>/etc/vant-agent/config.yaml</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/vant-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/vant-agent.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
